# build/click.service
[Unit]
Description=Click app

[Service]
WorkingDirectory=/var/www
ExecStart=/var/www/click run
Restart=on-failure
Type=simple

[Install]
Alias=click.service
WantedBy=multi-user.target
